# deb/etc/cron.d/snapraid-runner
0 1 * * 0   root    python3 /opt/snapraid-runner/snapraid-runner.py -c /etc/snapraid-runner.conf
