# deb/etc/cron.d/snapraid-runner
0 1 * * *   root    python2 /opt/snapraid-runner/snapraid-runner.py -c /etc/snapraid-runner.conf
