# deb/etc/cron.d/snapraid-runner
00 08 * * * python2 /opt/snapraid-runner/snapraid-runner.py -c /etc/snapraid-runner.conf